# com.fleming.ai.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Service identifier -->
    <key>Label</key>
    <string>com.fleming.ai</string>

    <!-- Program to run -->
    <key>ProgramArguments</key>
    <array>
        <string>/opt/anaconda3/bin/python3</string>
        <string>/Users/jnnj92/Fleming-AI/main.py</string>
    </array>

    <!-- Working directory -->
    <key>WorkingDirectory</key>
    <string>/Users/jnnj92/Fleming-AI</string>

    <!-- Start on load -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Keep alive - restart if crashes -->
    <key>KeepAlive</key>
    <dict>
        <!-- Only restart if process exits abnormally -->
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <!-- Restart throttle (seconds) -->
    <key>ThrottleInterval</key>
    <integer>60</integer>

    <!-- Standard output log -->
    <key>StandardOutPath</key>
    <string>/Users/jnnj92/Library/Logs/fleming-ai.log</string>

    <!-- Standard error log -->
    <key>StandardErrorPath</key>
    <string>/Users/jnnj92/Library/Logs/fleming-ai-error.log</string>

    <!-- Environment variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/anaconda3/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>

    <!-- Process type -->
    <key>ProcessType</key>
    <string>Background</string>

    <!-- Nice value (lower priority) -->
    <key>Nice</key>
    <integer>5</integer>

    <!-- Disable automatic termination -->
    <key>AbandonProcessGroup</key>
    <true/>
</dict>
</plist>
